# dots/.config/systemd/user/ii-lacuna-upstream-check.timer
[Unit]
Description=Weekly check for ii-lacuna upstream updates

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
