# dots/.config/systemd/user/ii-lacuna-upstream-check.service
[Unit]
Description=Check ii-lacuna upstream for updates

[Service]
Type=oneshot
ExecStart=/home/donnie/projects/ii-lacuna/scripts/check-upstream.sh
